# roles/fozzie_drives/templates/snapraid-btrfs-runner.timer.j2
[Unit]
Description=Run snapraid-btrfs-runner every night

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
