# roles/fozzie_drives/templates/snapraid-btrfs-runner.service.j2
[Unit]
Description=Run snapraid-btrfs-runner every night

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ snapraid_btrfs_runner_bin }} -c {{ snapraid_btrfs_runner_config }}
